# service-units/systemd/opensquilla.service
# OpenSquilla Linux systemd --user unit — user-scope only.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp service-units/systemd/opensquilla.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now opensquilla.service
#
# Safety contract — DO NOT REMOVE:
#   OPENSQUILLA_LISTEN=127.0.0.1 binds the gateway to loopback only. Network
#   exposure is an explicit operator opt-in. To expose the gateway on the
#   network, create a drop-in override at
#     ~/.config/systemd/user/opensquilla.service.d/override.conf
#   with:
#     [Service]
#     Environment="OPENSQUILLA_LISTEN=0.0.0.0"
#   then re-run `systemctl --user daemon-reload` + `systemctl --user restart
#   opensquilla.service`. Only run with OPENSQUILLA_LISTEN=0.0.0.0 behind a trusted
#   reverse proxy or VPN — opensquilla's first-class auth assumes loopback scope
#   by default. The S19 boot WARNING (gateway.bind.public) will fire whenever
#   the unit is configured for wildcard bind.

[Unit]
Description=OpenSquilla agent runtime gateway (loopback default)
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/opensquilla gateway run
Environment="OPENSQUILLA_LISTEN=127.0.0.1"
Environment="OPENSQUILLA_GATEWAY_PORT=18790"
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
